# system/hardware/service/usb-control/src/main/resources/usb-control.rules
ACTION=="add", SUBSYSTEM=="block", ENV{ID_BUS}=="usb", RUN+="/usr/mediaserver/sbin/usb-control add %k"
ACTION=="remove", SUBSYSTEM=="block", ENV{ID_BUS}=="usb", RUN+="/usr/mediaserver/sbin/usb-control remove %k"