# system/hardware/service/fan-control/src/main/resources/fan-control.service
[Unit]
Description=fan-control service
[Service]
ExecStart=/usr/local/sbin/fan-control
[Install]
WantedBy=multi-user.target